# helloworld.conf
description "Yeehaw!"

# Start when the system is ready to do networking.
start on started elastic-network-interfaces

# Stop when the system is on its way down.
stop on shutdown

respawn
script
    exec /usr/bin/node /home/ec2-user/helloworld.js
end script
